# systemd/gpu-eco-mode.service
[Unit]
Description=Set NVIDIA GPU to ECO Mode on Boot
After=nvidia-persistenced.service
Wants=nvidia-persistenced.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 -m gpu_thermal.cli eco
WorkingDirectory=/home/jeb/programs/gpu-thermal-worlock
Environment=PYTHONPATH=/home/jeb/programs/gpu-thermal-worlock
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
